# d03/setup/systemd/iscsid.service.d/network-online.conf
[Unit]
# Ensure iscsid starts after network is up so iSCSI login succeeds before mount-docker-iscsi runs
After=network-online.target
Wants=network-online.target
