# d01/setup/systemd/open-iscsi.service.d/override.conf
# Debian Trixie: open-iscsi 2.1.9+ stores nodes in /var/lib/iscsi/nodes but
# open-iscsi.service still checks /etc/iscsi/nodes, so the service never runs
# and automatic login at boot never happens (Debian bug #1090725).
# Override so the service runs when we have nodes in the new location.
[Unit]
ConditionDirectoryNotEmpty=
ConditionDirectoryNotEmpty=/var/lib/iscsi/nodes
